# sys/systemd/kanban.service
[Unit]
Description=Kanban Board Service
After=network.target

[Service]
Type=simple
User=kanban
Group=kanban
WorkingDirectory=/opt/kanban
Environment=PATH=/opt/kanban/venv/bin
Environment=DATABASE_PATH=/opt/kanban/kanban.db
Environment=STATIC_PATH=/opt/kanban/backend/static
ExecStart=/opt/kanban/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --app-dir /opt/kanban
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kanban

[Install]
WantedBy=multi-user.target